# deploy/lobsec-backup.service
[Unit]
Description=lobsec daily backup
After=lobsec.service

[Service]
Type=oneshot
User=lobsec
Group=lobsec
Environment=SOFTHSM2_CONF=/opt/lobsec/boot/softhsm2.conf
ExecStart=/opt/lobsec/bin/backup.sh
