# deploy/lobsec-examy-test.timer
[Unit]
Description=Examy QA test daily execution timer (3am UTC)

[Timer]
OnCalendar=*-*-* 03:00:00
AccuracySec=1min
Persistent=false
RandomizedDelaySec=0

[Install]
WantedBy=timers.target
