# deploy/lobsec-watchdog.service
[Unit]
Description=lobsec stack watchdog — restart dead services
Documentation=https://github.com/BlackDuck314/lobsec

[Service]
Type=oneshot
ExecStart=/opt/lobsec/bin/watchdog.sh
